# services/orchestrator/orchestrator.service
[Unit]
Description=OpenClaw Orchestrator Service
After=network.target postgresql.service
Documentation=https://github.com/Ashwin012/openclaw-dashboard

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/projects/dev-dashboard/services/orchestrator
Environment=NODE_ENV=production
EnvironmentFile=-/home/openclaw/projects/dev-dashboard/.env
ExecStart=/home/openclaw/.nvm/versions/node/v20.20.1/bin/node server.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orchestrator

[Install]
WantedBy=multi-user.target
